# resources/fty-metric-tpower.service.in
[Unit]
Description=Total power agent for 42ITy project
After=malamute.service fty-db.target
Requires=malamute.service fty-db.target
PartOf=bios.target

[Service]
Type=simple
User=@AGENT_USER@
Restart=always
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/usr/share/fty/etc/default/fty
EnvironmentFile=-/usr/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/fty
EnvironmentFile=-/etc/default/fty__%n.conf
Environment="prefix=/usr"
EnvironmentFile=/etc/default/bios-db-ro
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/@PROJECT_NAME@ @AGENT_CONF_FILE@
Restart=always

[Install]
WantedBy=bios.target
