# src/fty-metric-tpower.service.in
[Unit]
Description=Total power agent for 42ITy project
After=malamute.service bios-db-init.service
Requires=malamute.service bios-db-init.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Environment="prefix=@prefix@"
EnvironmentFile=-/etc/default/bios-db-ro
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/fty__%n.conf
ExecStart=@prefix@/bin/fty-metric-tpower @sysconfdir@/@PACKAGE@/fty-metric-tpower.cfg
Restart=always

[Install]
WantedBy=bios.target
