# src/fty-metric-tpower.service.in
[Unit]
Description=fty-metric-tpower service
After=network.target

[Service]
Type=simple
Environment="prefix=@prefix@"
EnvironmentFile=-/etc/default/bios-db-ro
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/fty__%n.conf
ExecStart=@prefix@/bin/fty-metric-tpower @sysconfdir@/@PACKAGE@/fty-metric-tpower.cfg
Restart=always

[Install]
WantedBy=bios.target
