# src/fty-metric-tpower.service.in
[Unit]
Description=Total power agent for 42ITy project
After=malamute.service bios-db-init.service
Requires=malamute.service bios-db-init.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
EnvironmentFile=-@sysconfdir@/default/bios-db-ro
ExecStart=@prefix@/bin/fty-metric-tpower @sysconfdir@/@PACKAGE@/fty-metric-tpower.cfg
Restart=always

[Install]
WantedBy=bios.target
